# synthos_build/ops/systemd/pi5/synthos-portal.service
[Unit]
Description=Synthos Retail Portal (Gunicorn — port 5001)
# Ordered to start AFTER the boot-sequence oneshot, so env / DB / health
# checks have been exercised before the portal accepts traffic.  Wants=
# is a soft dep: if boot_sequence fails the portal still starts (matches
# the pre-systemd philosophy of "don't block agents; let them fail
# gracefully and log their own errors").
After=network.target synthos-boot-sequence.service
Wants=network.target synthos-boot-sequence.service

[Service]
Type=simple
User=pi516gb
Group=pi516gb
WorkingDirectory=/home/pi516gb/synthos/synthos_build/src
ExecStart=/home/pi516gb/.local/bin/gunicorn \
    --workers 2 \
    --threads 2 \
    --bind 0.0.0.0:5001 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    retail_portal:app
# Bind 0.0.0.0 is intentional — pi4b's cloudflared has
#   portal.synth-cloud.com → http://10.0.0.11:5001
# in /etc/cloudflared/config.yml, so the public ingress crosses the LAN.
# Tried 127.0.0.1 binding 2026-04-24 and broke the portal — see audit
# Phase 5 + the recommended architectural fix (move portal hostname to
# pi5's own tunnel, which already terminates app.synth-cloud.com via
# localhost, and the bind tightening becomes safe).
Restart=always
RestartSec=5
StandardOutput=append:/home/pi516gb/synthos/synthos_build/logs/portal.log
StandardError=append:/home/pi516gb/synthos/synthos_build/logs/portal.log

# Env file loaded by the portal itself via python-dotenv
# Do not duplicate secrets here

[Install]
WantedBy=multi-user.target
