# synthos_build/ops/systemd/pi5/synthos-market-daemon.service
[Unit]
Description=Synthos Retail Market Daemon (enrichment pipeline, 9:10-16:00 ET)
# Kicked by synthos-market-daemon.timer each weekday at 09:10 ET. The daemon
# runs the enrichment pipeline (screener → news → sentiment → macro →
# market_state → bias → validator → price_poller) on a 30-min cadence and
# exits cleanly at 16:00 ET. Trader dispatch is handled by
# synthos-trade-daemon (separate process).
After=network.target synthos-boot-sequence.service
Wants=network.target synthos-boot-sequence.service

[Service]
Type=simple
User=pi516gb
Group=pi516gb
WorkingDirectory=/home/pi516gb/synthos/synthos_build
ExecStart=/usr/bin/python3 src/retail_market_daemon.py
# Restart only on non-zero exit (crash). Normal 16:00 exit = 0, systemd does
# NOT restart — the timer re-triggers tomorrow at 09:10.
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/pi516gb/synthos/synthos_build/logs/market_daemon.log
StandardError=append:/home/pi516gb/synthos/synthos_build/logs/market_daemon.log

[Install]
# NOTE: do NOT enable this service directly — enable the .timer instead.
# multi-user.target listed only so `systemctl status` has a proper target.
WantedBy=multi-user.target
