# synthos_build/ops/systemd/pi4b/synthos-company-server.service
[Unit]
Description=Synthos Company Admin Server
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/synthos-company
EnvironmentFile=/home/pi/synthos-company/company.env
ExecStart=/usr/bin/python3 /home/pi/synthos-company/company_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
