# synthos_build/config/systemd/synthos-watchdog.service
[Unit]
Description=Synthos Retail Watchdog (crash monitor + recovery)
After=network.target synthos-portal.service
Wants=synthos-portal.service

[Service]
Type=simple
User=pi516gb
Group=pi516gb
WorkingDirectory=/home/pi516gb/synthos/synthos_build/src
ExecStart=/usr/bin/python3 retail_watchdog.py
Restart=always
RestartSec=10
StandardOutput=append:/home/pi516gb/synthos/synthos_build/logs/watchdog.log
StandardError=append:/home/pi516gb/synthos/synthos_build/logs/watchdog.log

[Install]
WantedBy=multi-user.target
