# synthos_build/config/systemd/synthos-watchdog.service
[Unit]
Description=Synthos Retail Watchdog (crash monitor + recovery)
After=network.target synthos-portal.service
Wants=synthos-portal.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/synthos/synthos_build/src
ExecStart=/home/pi/synthos/venv/bin/python3 retail_watchdog.py
Restart=always
RestartSec=10
StandardOutput=append:/home/pi/synthos/synthos_build/logs/watchdog.log
StandardError=append:/home/pi/synthos/synthos_build/logs/watchdog.log

[Install]
WantedBy=multi-user.target
