# synthos_build/config/systemd/synthos-scheduler@.service
[Unit]
Description=Synthos Scheduler — session %i
After=network.target

[Service]
Type=oneshot
User=pi516gb
Group=pi516gb
WorkingDirectory=/home/pi516gb/synthos/synthos_build/src
ExecStart=/usr/bin/python3 retail_scheduler.py --session %i
StandardOutput=append:/home/pi516gb/synthos/synthos_build/logs/scheduler.log
StandardError=append:/home/pi516gb/synthos/synthos_build/logs/scheduler.log

# oneshot services exit when done — no Restart needed.
# systemd records success/failure in the journal automatically.
