# synthos_build/config/systemd/synthos-portal.service
[Unit]
Description=Synthos Retail Portal (Flask — port 5001)
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi516gb
Group=pi516gb
WorkingDirectory=/home/pi516gb/synthos/synthos_build/src
ExecStart=/usr/bin/python3 retail_portal.py
Restart=always
RestartSec=5
StandardOutput=append:/home/pi516gb/synthos/synthos_build/logs/portal.log
StandardError=append:/home/pi516gb/synthos/synthos_build/logs/portal.log

# Env file loaded by the portal itself via python-dotenv
# Do not duplicate secrets here

[Install]
WantedBy=multi-user.target
